# deploy/systemd/openclaw-gateway.service.example
[Unit]
Description=OpenClaw Gateway Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=<YOUR_LINUX_USER>
WorkingDirectory=/home/<YOUR_LINUX_USER>
ExecStart=/usr/bin/openclaw gateway start
ExecStop=/usr/bin/openclaw gateway stop
ExecReload=/usr/bin/openclaw gateway restart
Restart=always
RestartSec=5

# Optional hardening (enable if compatible with your setup)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=false

[Install]
WantedBy=multi-user.target
